# hw17-backup/provision/roles/borgclient-role/files/backup.service
[Unit]
Description=My backup service

[Service]
Type=oneshot
#EnvironmentFile=/etc/sysconfig/watchlog.conf
ExecStart=/opt/backup.sh
